# deploy/ai_image_analyzer.service.example
[Unit]
Description=AI Image Analyzer bot (systemd service example)
After=network.target

[Service]
# Replace 'aiimage' with the system user that will run the bot
User=aiimage
Group=aiimage

# Set WorkingDirectory to the repository path on the server
WorkingDirectory=/path/to/ai_image_analyzer

# Load environment variables from the .env file (optional)
EnvironmentFile=/path/to/ai_image_analyzer/.env

# ExecStart should point to your python venv interpreter and bot.py
ExecStart=/path/to/venv/bin/python /path/to/ai_image_analyzer/bot.py

Restart=on-failure
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=20

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

# Usage:
# 1) copy this file to /etc/systemd/system/ai_image_analyzer.service
# 2) create system user (example): sudo useradd -m aiimage
# 3) set ownership of the repo to that user, create venv and install requirements
# 4) sudo systemctl daemon-reload
# 5) sudo systemctl enable --now ai_image_analyzer.service